# deploy/clearcontract-web.service
[Unit]
Description=ClearContract Next.js web
After=network.target pocketbase.service

[Service]
Type=simple
User=clearcontract
WorkingDirectory=/opt/clearcontract/web
EnvironmentFile=/etc/clearcontract/web.env
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
